# systemd/firelight-rest.service
[Unit]
Description=Daemon for LED strip control via ws2811

[Service]
ExecStart=/usr/bin/firelight-rest

[Install]
WantedBy=multi-user.target

